# build/rootfs/etc/services.d/99-homeassistant/run
#!/usr/bin/with-contenv bash
s6-setuidgid abc
cd /usr/src/app

sleep 10
s6-svwait -t 5000 -U /var/run/s6/services/socat

exec python -m homeassistant --config /config
